# docker/root/etc/s6-overlay/s6-rc.d/svc-cron/run
#!/command/execlineb -P
crond -f -L /dev/stderr
